# systemd/user/xscreensaver.service
[Unit]
Description=XscreenSaver
[Service]
ExecStart=/usr/bin/xscreensaver -no-splash
[Install]
WantedBy=default.target
